# docker/pycron/script/crontabs
59 23 * * * flock -xn /tmp/worldometers_gdp.lock -c 'timeout 5m /usr/bin/python3 /data/script/worldometers/gdp.py >> /data/script/worldometers/log/gdp.log 2>&1'
